# app/rcstation/scripts/rcstation.target
[Unit]
Description=Make the comm system of the RC station operational
AllowIsolate=no
