# app/rcstation/scripts/rcstation_udp_rx_video.service
[Unit]
Description=RC Car Video Receiver
Requires= rcstation_configure_wifi.service
After= rcstation_configure_wifi.service
PartOf=rcstation.target 

[Service]
ExecStart=/usr/local/bin/rcstation_udp_rx_video.bash
EnvironmentFile=/etc/rccar.conf

[Install]
WantedBy=rcstation.target 
