# app/rccar/scripts/rccar.target
[Unit]
Description=Make the RC car operational
Requires=multi-user.target
After=multi-user.target
AllowIsolate=yes
